# etc/systemd/system/gs.snap.lobby.sync-SnapshotVR@.path
[Unit]
Description=%j #%i Lobby Sync Path


[Install]
WantedBy=multi-user.target


[Path]
PathExistsGlob=/opt/snap-gs/%j/%i/log/*.gz
