# etc/systemd/system/gs.snap.lobby.idle-SnapshotVR@.path
[Unit]
Description=%j #%i Lobby Idle Path


[Install]
WantedBy=multi-user.target


[Path]
PathExists=/opt/snap-gs/%j/%i/stat/idle
PathChanged=/opt/snap-gs/%j/%i/stat/idle
PathModified=/opt/snap-gs/%j/%i/stat/idle
PathExistsGlob=/opt/snap-gs/%j/%i/stat/idle
