# etc/systemd/system/gs.snap.lobby-SnapshotVR@.path
[Unit]
Description=%j #%i Lobby Path


[Install]
WantedBy=multi-user.target


[Path]
PathChanged=/opt/snap-gs/%j/%i/flags/start
PathChanged=/opt/snap-gs/%j/%i/spec/start

PathChanged=/opt/snap-gs/%j/%i/spec/up
PathExists=/opt/snap-gs/%j/%i/spec/up

PathChanged=/opt/snap-gs/%j/%i/peer/full
PathExists=/opt/snap-gs/%j/%i/peer/full
