# etc/systemd/system/gs.snap.lobby-SnapshotVR@.path
[Unit]
Description=%j #%i Lobby Path


[Install]
WantedBy=multi-user.target


[Path]
PathExists=/opt/snap-gs/%j/%i/peer/full
PathChanged=/opt/snap-gs/%j/%i/peer/full
PathModified=/opt/snap-gs/%j/%i/peer/full
PathExistsGlob=/opt/snap-gs/%j/%i/peer/full
